# goathub.service
# GoatHub Trading Bot — systemd Unit.
#
# 2026-06-12 Hardening (#8 + Audit M-16/M-19/LOW-11): läuft als goathub statt
# root, bindet nur noch 127.0.0.1 (Caddy proxied + terminiert TLS), systemd-
# Sandbox. Bei RCE/kompromittierter Dependency hängt sonst .env (JWT_SECRET,
# ENCRYPTION_KEY, DISCORD_BOT_TOKEN) + goathub.db (alle Agent-Keys) an root.
#
# ── EINMALIGE Server-Voraussetzungen (als root, BEVOR dieses Unit live geht) ──
#   useradd --system --no-create-home --shell /usr/sbin/nologin goathub
#   mkdir -p /var/lib/goathub && chown goathub:goathub /var/lib/goathub && chmod 700 /var/lib/goathub
#   mkdir -p /var/backups/goathub && chown goathub:goathub /var/backups/goathub && chmod 700 /var/backups/goathub
#   chown -R goathub:goathub /var/www/goathub-bot
#   git config --global --add safe.directory /var/www/goathub-bot   # root pullt im goathub-Repo (deploy.yml setzt das auch selbst)
#   cp goathub.service goathub-backup.service goathub-backup.timer /etc/systemd/system/
#   systemctl daemon-reload && systemctl restart goathub && systemctl enable --now goathub-backup.timer
# Danach installiert der Deploy-Workflow Unit-Updates automatisch (gated auf
# `id goathub` — solange der User fehlt, bleibt das alte Unit aktiv).
#
# Env-Handling: KEIN EnvironmentFile= — app/config.py lädt die .env selbst via
# python-dotenv aus dem WorkingDirectory (/var/www/goathub-bot/.env).
[Unit]
Description=GoatHub Trading Bot (Multi-User Plattform)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=goathub
Group=goathub
WorkingDirectory=/var/www/goathub-bot
# Bind auf 127.0.0.1 statt 0.0.0.0 — Caddy proxied lokal auf :8000
# (TLS-Terminierung). Direkter Internet-Zugriff auf uvicorn umging vorher
# Caddys X-Forwarded-For → Rate-Limit-Key in main.py wäre spoofbar gewesen
# (Pair-Fix mit _client_ip in app/main.py).
ExecStart=/var/www/goathub-bot/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
# 2026-06-13 Review-Fix: zurück auf always. `systemctl stop` restartet bei
# systemd NIE automatisch (Restart= gilt nur für unrequested exits) — der
# Operator-Halt war also nie das Problem. on-failure hätte dagegen einen
# unerwarteten sauberen Selbst-Exit (exit 0, z.B. stray SIGTERM mit graceful
# uvicorn-Shutdown) liegen lassen: Engine down, kein Sync, niemand merkt's.
# Für eine Trading-Engine ist Verfügbarkeit der sichere Default.
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# systemd-Sandbox.
NoNewPrivileges=true
# ProtectSystem=strict mountet das gesamte FS read-only (strenger als das im
# Audit geforderte "full"); beschreibbar bleiben nur App-Dir (Code, venv,
# goathub.db inkl. SQLite-WAL/-Journal) + /var/lib/goathub (Emergency-Halt-
# Flag, LOW-11) + das private /tmp.
ProtectSystem=strict
ReadWritePaths=/var/www/goathub-bot /var/lib/goathub
ProtectHome=true
# 2026-06-12 LOW-11: PrivateTmp jetzt true — das EMERGENCY_HALT_FLAG liegt
# nicht mehr unter world-writable /tmp sondern unter /var/lib/goathub/
# emergency-halt (config.py). Manueller Not-Halt per Shell entsprechend:
#   touch /var/lib/goathub/emergency-halt
# VORAUSSETZUNG: /var/lib/goathub existiert + gehört goathub (siehe oben),
# sonst fällt config.py auf /tmp zurück und der manuelle Halt-Pfad wäre durch
# PrivateTmp unsichtbar.
PrivateTmp=true

[Install]
WantedBy=multi-user.target
